# v3/v3-offline.service
[Unit]
Description=MisterMoney V3 Offline Worker
After=network.target postgresql.service redis.service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/.openclaw/workspace/MisterMoney
ExecStart=/usr/bin/python3 -m v3.offline.main
Restart=always
RestartSec=30
Environment=PYTHONPATH=/home/ubuntu/.openclaw/workspace/MisterMoney

[Install]
WantedBy=default.target
